# deployment/puppet/nailgun/templates/systemd/nailgun.service.erb
[Unit]
Name=Nailgun daemon
ConditionPathExists=/etc/nailgun/uwsgi_nailgun.yaml

[Service]
Restart=on-failure
RestartSec=10s
StartLimitBurst=10
TimeoutStartSec=10s
TimeoutStopSec=5s

StandardInput=null
StandardOutput=syslog
StandardError=syslog+console
SyslogIdentifier=nailgun

ExecStart=/usr/sbin/uwsgi -y /etc/nailgun/uwsgi_nailgun.yaml
ExecReload=/usr/sbin/uwsgi --reload /var/run/nailgun.pid
ExecStop=/usr/sbin/uwsgi --stop /var/run/nailgun.pid
ExecStopPost=/usr/bin/rm -f /var/run/nailgun.pid

[Install]
WantedBy=multi-user.target

