# deployment/puppet/nailgun/templates/systemd/receiverd.service.erb
[Unit]
Name=Receiverd daemon

[Service]
Restart=on-failure
RestartSec=10s
StartLimitBurst=10
TimeoutStartSec=10s
TimeoutStopSec=5s

StandardInput=null
StandardOutput=syslog
StandardError=syslog+console
SyslogIdentifier=receiverd

ExecStart=/usr/bin/receiverd

[Install]
WantedBy=multi-user.target

