# deployment/puppet/nailgun/templates/systemd/statsenderd.service.erb
[Unit]
Name=Statsenderd daemon

[Service]
Restart=on-failure
RestartSec=10s
StartLimitBurst=10
TimeoutStartSec=10s
TimeoutStopSec=5s

StandardInput=null
StandardOutput=syslog
StandardError=syslog+console
SyslogIdentifier=statsenderd

ExecStart=/usr/bin/statsenderd

[Install]
WantedBy=multi-user.target

