# deployment/puppet/docker/templates/systemd/template.service.erb
[Unit]
Name=<%= @title %> container
Requires=docker.service
After=docker.service <% if @depends[@title] -%>docker-<%= @depends[@title] -%>.service<% end %>

[Service]
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitInterval=60
ExecStartPre=/usr/bin/dockerctl create <%= @title %>
ExecStartPre=/usr/bin/docker start fuel-core-<%= @release %>-<%= @title %>
ExecStartPre=/usr/bin/dockerctl check <%= @title %>
ExecStart=/usr/bin/docker attach --no-stdin=true fuel-core-<%= @release %>-<%= @title %>
ExecStop=/usr/bin/docker stop -t <%= @timeout %> fuel-core-<%= @release %>-<%= @title %>

[Install]
WantedBy=multi-user.target
