# deployment/puppet/cluster/files/rabbit-fence.upstart
description "RabbitMQ-fence daemon"

start on runlevel [2345]
stop on runlevel [016]

exec /usr/bin/rabbit-fence.py
