# debian/fuel-rabbit-fence.upstart
description "RabbitMQ-fence daemon"

start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 10 5

pre-start script
  mkdir -p /var/run/rabbitmq
  chown -R rabbitmq /var/run/rabbitmq

  # TODO(dmburmistrov): remove this after fixing PID-file handling
  rm -f /var/run/rabbitmq/rabbit-fence.pid
end script

exec /usr/bin/rabbit-fence.py
