# debian/fuel-umm.umm-tr.init.in
#! /bin/sh
### BEGIN INIT INFO
# Provides:          umm-tr
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: umm-tr
# Description:       umm-tr
### END INIT INFO

# Author: Thomas Goirand <zigo@debian.org>
DESC="Openstack maintenance mode - time reset"
PROJECT_NAME=fuel
NAME=umm-tr
NO_OPENSTACK_CONFIG_FILE_DAEMON_ARG=yes
USE_LOGFILE=no
SYSTEMD_TYPE=oneshot
SYSTEM_USER=root
SYSTEM_GROUP=root
