# deployment/puppet/cgroups/templates/cgrulesengd.conf.erb
description "CGroups Rules Engine Daemon"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork

pre-start script
test -x /usr/sbin/cgrulesengd || { stop; exit 0; }
end script

script

# get default options
if [ -f "/etc/default/cgrulesengd" ] ; then
   . /etc/default/cgrulesengd
   OPTIONS="$NODAEMON $LOG"
   if [ -n "$LOG_FILE" ]; then
      OPTIONS="$OPTIONS --logfile=$LOG_FILE"
   fi
   if [ -n "$SOCKET_USER" ]; then
      OPTIONS="$OPTIONS -u $SOCKET_USER"
   fi
   if [ -n "$SOCKET_GROUP" ]; then
      OPTIONS="$OPTIONS -g $SOCKET_GROUP"
   fi
else
   OPTIONS=""
fi

# Don't run if no configuration file
CGRED_CONF=/etc/cgrules.conf
if [ ! -s "$CGRED_CONF" ]; then
   echo "CGroups Rules Engine Daemon unconfigured"
   stop
   exit 0
fi

#exec /usr/sbin/cgrulesengd --nodaemons $OPTIONS
exec cgrulesengd $OPTIONS

end script
