# files/fuel-umm/root/etc/profile.d/umm.sh
if [ "z`umm status`" = "zumm" ] ; then
    cat /etc/issue.mm
    [ "$#" != "0" ] && echo "$@"
fi
